# nixos/HOME/1/HOME,.config,systemd,user,plasma-kwin_wayland.service
[Unit]
Description=KDE Window Manager
PartOf=graphical-session.target

[Service]
ExecStart=/run/current-system/sw/bin/wayfire
Slice=session.slice
